# supervisord.conf
[program:RouterDriver]
command=/usr/bin/gunicorn -b 0.0.0.0:80 main:app
directory=/root/physical_router_driver
user=root
autostart=true
